# supervisord.d/uwsgi.ini
[program:uwsgi]
user = uwsgi
command=/data/miya/config/start_uwsg.sh
autostart=true
autorestart=true
stderr_logfile = /data/logs/uwsgi.error.log
stdout_logfile = /data/logs/uwsgi.log
stopsignal=INT
